# meta-resin-common/recipes-support/supervisor-init/files/supervisor-init.service
[Unit]
Description=Resin supervisor service
After=connman.service

[Service]
Type=forking
ExecStart=-@BASE_BINDIR@/supervisor-init
StandardError=syslog
RemainAfterExit=no

[Install]
WantedBy=basic.target
