# meta-resin-common/recipes-support/supervisor-init/files/supervisor-init.service
[Unit]
Description=Resin supervisor service
Wants=connman.service resin-init.service
After=connman.service resin-init.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=-@BASE_BINDIR@/supervisor-init
StandardError=journal+console

[Install]
WantedBy=basic.target
