# meta-balena-common/recipes-support/development-features/development-features/development-features.target
[Unit]
Description=Development features

[Install]
WantedBy=multi-user.target
