# meta-balena-common/recipes-support/balena-persistent-logs/balena-persistent-logs/balena-persistent-logs.service
[Unit]
Description=Balena persistent logs
Requires=systemd-journald.service
After=systemd-journald.service
RequiresMountsFor=/mnt/state /mnt/boot

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@BINDIR@/balena-persistent-logs

[Install]
WantedBy=multi-user.target
