# meta-balena-common/recipes-core/systemd/systemd/10-zram.rules
KERNEL=="zram0", SUBSYSTEM=="block", ACTION=="add", RUN+="/lib/udev/zram-swap-init $env{DEVNAME}"
