# meta-balena-common/recipes-support/balena-config-vars/balena-config-vars/config-json.path
[Unit]
Description=Config.json path watch

[Path]
PathChanged=/mnt/boot/config.json

[Install]
WantedBy=basic.target
