# meta-balena-common/recipes-support/balena-config-vars/balena-config-vars/config-json.target
[Unit]
Description=Config.json watcher target
