# meta-balena-common/recipes-support/balena-config-vars/balena-config-vars/config-json.service
[Unit]
Description=Config.json watcher service

[Service]
Type=oneshot
ExecStart=/bin/echo 'config.json changed'
