# meta-balena-common/recipes-support/balena-config-vars/balena-config-vars/os-udevrules.service
[Unit]
Description=udev rules configuration from config.json
Requires=resin-boot.service bind-etc-udev-rules.d.service
After=resin-boot.service bind-etc-udev-rules.d.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@SBINDIR@/os-udevrules

[Install]
WantedBy=multi-user.target
