# meta-resin-common/recipes-support/resin-supervisor/files/update-resin-supervisor.timer
[Unit]
Description=Resin supervisor updater timer

[Timer]
OnBootSec=15min
OnUnitInactiveSec=5min

[Install]
WantedBy=multi-user.target
