# meta-resin-common/recipes-support/resin-mounts/resin-mounts-flasher/bind-path.service.in
[Unit]
Description=Flasher bind mount for @target@
DefaultDependencies=no
Requires=tmp.mount var-volatile-lib.service systemd-tmpfiles-setup.service
After=tmp.mount var-volatile-lib.service systemd-tmpfiles-setup.service
Before=umount.target
Conflicts=umount.target

[Service]
Type=oneshot
RemainAfterExit=Yes
ExecStart=/usr/bin/bindmount --target @target@ --bind-root /tmp/root-overlay --command mount
ExecStop=/usr/bin/bindmount --target @target@ --bind-root /tmp/root-overlay --command unmount

[Install]
WantedBy=multi-user.target
