# meta-resin-common/recipes-support/resin-mounts/resin-mounts-flasher/etc-openvpn-resin.conf.mount
[Unit]
Description=Resin openvpn configuration
Requires=prepare-openvpn.service
After=prepare-openvpn.service

[Mount]
What=/run/openvpn/resin.conf
Where=/etc/openvpn/resin.conf
Type=none
Options=bind

[Install]
WantedBy=multi-user.target
