# meta-resin-common/recipes-support/resin-mounts/resin-mounts-flasher/resin-boot.service
[Unit]
Description=Resin boot partition mount compatibility service
DefaultDependencies=no
Requires=mnt-boot.mount temp-conf.service
After=mnt-boot.mount temp-conf.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/echo "Boot partition mounted"
ExecStop=/bin/echo "This is a service wrapper. Use mnt-boot.mount instead."

[Install]
WantedBy=multi-user.target
