# meta-balena-common/recipes-connectivity/balena-proxy-config/balena-proxy-config/balena-proxy-config.service
[Unit]
Description=Resin proxy configuration service
After=dnsmasq.service

[Service]
Type=oneshot
ExecStart=@BASE_BINDIR@/sh @BINDIR@/balena-proxy-config
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
