# meta-balena-common/recipes-connectivity/balena-proxy-config/balena-proxy-config/redsocks-conf.target
[Unit]
Description=redsocks.conf watcher target
