# meta-balena-common/recipes-connectivity/bluez5/files/10-local-bt-hci-up.rules
ACTION=="add", SUBSYSTEM=="bluetooth", KERNEL=="hci[0-9]*", RUN+="/usr/bin/hciconfig %k up"
