# meta-balena-common/recipes-connectivity/balena-net-config/balena-net-config/balena-net-config.service
[Unit]
Description=Resin network configure service
Requires=resin-boot.service bind-etc-NetworkManager-system-connections.service
After=resin-boot.service bind-etc-NetworkManager-system-connections.service
PartOf=config-json.target

[Service]
ExecStart=@BASE_BINDIR@/bash @BINDIR@/balena-net-config
Type=oneshot
RemainAfterExit=yes
StandardError=journal+console

[Install]
WantedBy=multi-user.target
