# meta-balena-common/recipes-support/resin-mounts/resin-mounts/mnt-sysroot-active.automount
[Unit]
Description=mnt-sysroot-active mount
ConditionVirtualization=!docker

[Automount]
Where=/mnt/sysroot/active

[Install]
WantedBy=multi-user.target
