# meta-balena-common/recipes-support/resin-mounts/resin-mounts/mnt-sysroot-inactive.automount
[Unit]
Description=mnt-sysroot-inactive partition
ConditionVirtualization=!docker

[Automount]
Where=/mnt/sysroot/inactive
