# meta-balena-common/recipes-support/resin-mounts/resin-mounts/mnt-data.mount
[Unit]
Description=mnt-data mount
ConditionVirtualization=!docker

[Mount]
What=/dev/disk/by-state/resin-data
Where=/mnt/data
Options=noauto
