# meta-balena-common/recipes-support/resin-mounts/resin-mounts/mnt-data.automount
[Unit]
Description=mnt-data automount
ConditionVirtualization=!docker

[Automount]
Where=/mnt/data

[Install]
WantedBy=multi-user.target
