# meta-balena-common/recipes-support/resin-mounts/resin-mounts/mnt-state.automount
[Unit]
Description=mnt-state automount
ConditionVirtualization=!docker

[Automount]
Where=/mnt/state
