# meta-balena-common/recipes-support/resin-mounts/resin-mounts/mnt-boot.automount
[Unit]
Description=mnt-boot automount
ConditionVirtualization=!docker

[Automount]
Where=/mnt/boot
