# meta-balena-common/recipes-support/resin-mounts/resin-mounts/mnt-state.mount
[Unit]
Description=mnt-state mount
ConditionVirtualization=!docker

[Mount]
What=/dev/disk/by-state/resin-state
Where=/mnt/state
Options=noauto
