# meta-balena-common/recipes-support/resin-mounts/resin-mounts/mnt-sysroot-active.mount
[Unit]
Description=mnt-sysroot-active mount
ConditionVirtualization=!docker

[Mount]
What=/dev/disk/by-state/active
Where=/mnt/sysroot/active
Options=noauto
