# meta-balena-common/recipes-support/resin-mounts/resin-mounts/mnt-efi.automount
[Unit]
Description=mnt-efi automount
ConditionVirtualization=!docker

[Automount]
Where=/mnt/efi
