# meta-balena-common/recipes-support/resin-vars/resin-vars/config-json.path
[Unit]
Description=Config.json path watch
Requires=resin-boot.service
After=resin-boot.service

[Path]
PathChanged=/mnt/boot/config.json

[Install]
WantedBy=basic.target
