# meta-resin-common/recipes-support/resin-mounts/resin-mounts/mnt-sysroot-inactive.mount
[Unit]
Description=mnt-sysroot-inactive partition

[Mount]
What=/dev/disk/by-state/inactive
Where=/mnt/sysroot/inactive
Type=ext4
Options=noauto
