# meta-resin-common/recipes-support/resin-mounts/resin-mounts/mnt-boot.mount
[Unit]
Description=Resin boot partition mountpoint

[Mount]
What=/dev/disk/by-label/resin-boot
Where=/mnt/boot
Type=vfat

[Install]
WantedBy=multi-user.target
