# meta-resin-common/recipes-support/resin-mounts/resin-mounts/mnt-sysroot-inactive.automount
[Unit]
Description=mnt-sysroot-inactive partition

[Automount]
Where=/mnt/sysroot/inactive

[Install]
WantedBy=multi-user.target
