# meta-resin-common/recipes-support/resin-mounts/resin-mounts/mnt-conf-systemd.automount
[Unit]
Description=Automount systemd presistent state

[Automount]
Where=/mnt/conf/systemd

[Install]
WantedBy=multi-user.target
