# meta-resin-common/recipes-support/resin-mounts/resin-mounts/resin-bind.target
[Unit]
Description=Resin bind target
Documentation=man:systemd.target(5)
