# meta-resin-common/recipes-support/resin-vars/resin-vars/config-json.service
[Unit]
Description=Config.json watcher service
Wants=prepare-openvpn.service

[Service]
Type=oneshot
ExecStart=/bin/echo 'config.json changed'
