# meta-resin-common/recipes-support/resin-vars/resin-vars/config-json.path
[Unit]
Description=Config.json path watch
Requires=mnt-boot.mount
After=mnt-boot.mount

[Path]
PathChanged=/mnt/boot/config.json

[Install]
WantedBy=mnt-boot.mount
