# meta-resin-common/recipes-support/resin-btrfs/resin-btrfs-balance/resin-btrfs-balance.service
[Unit]
Description=Resin BTRFS balance
Requires=mnt-data.mount
After=mnt-data.mount

[Service]
Type=simple
ExecStart=@BASE_BINDIR@/bash @BINDIR@/resin-btrfs-balance
