# meta-resin-common/recipes-support/resin-monitor/resin-monitor/resin-monitor.timer
[Unit]
Description=Resin monitor timer

[Timer]
[Timer]
OnBootSec=5min
OnUnitInactiveSec=2min

[Install]
WantedBy=multi-user.target
