# meta-resin-common/recipes-support/resin-monitor/resin-monitor/resin-monitor.service
[Unit]
Description=Resin monitor

[Service]
Type=simple
ExecStart=@BASE_BINDIR@/bash @BINDIR@/resin-monitor
