# meta-balena-common/recipes-connectivity/openssh/recover-sshd-socket/recover-sshd-socket.timer
[Unit]
Description=Periodically check and recover sshd.socket

[Timer]
OnCalendar=*:0/15
Persistent=true

[Install]
WantedBy=multi-user.target
