# meta-balena-common/recipes-core/systemd/periodic-vacuum-logs/periodic-vacuum-logs.timer
[Unit]
Description=Periodic vacuum of journald logs

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=multi-user.target
