# meta-balena-common/recipes-core/systemd/periodic-vacuum-logs/periodic-vacuum-logs.service
[Unit]
After=systemd-journald.service

[Service]
ExecStart=/bin/journalctl --vacuum-size=32M
