# meta-balena-common/recipes-containers/resin-supervisor/resin-supervisor/resin-data.mount
[Unit]
Description=resin-data bind mount
Requires=resin-data.service
After=resin-data.service
DefaultDependencies=no

[Mount]
What=/mnt/data/resin-data
Where=/resin-data
Type=none
Options=bind

[Install]
WantedBy=multi-user.target
