# meta-resin-common/recipes-containers/docker-disk/docker-resin-supervisor-disk/resin.target
[Unit]
Description=Resin Target
Documentation=man:systemd.target(5)
