# meta-resin-common/recipes-connectivity/openvpn/openvpn/prepare-openvpn.service
[Unit]
Description=Prepare OpenVPN Application On Resin
Requires=mnt-boot.mount resin-device-uuid.service resin-device-api-key.service var-volatile.mount
After=mnt-boot.mount resin-device-uuid.service resin-device-api-key.service var-volatile.mount

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@BINDIR@/prepare-openvpn

[Install]
WantedBy=config-json.service
