# meta-balena-common/recipes-core/resin-extra-udev-rules/files/79-wlan-power.rules
ACTION=="add", KERNEL=="wl*", SUBSYSTEM=="net", RUN+="/usr/sbin/iw dev $name set power_save off"
