# meta-balena-common/recipes-connectivity/resin-proxy-config/resin-proxy-config/resin-proxy-config.service
[Unit]
Description=Resin proxy configuration service
Requires=resin-boot.service
After=resin-boot.service dnsmasq.service

[Service]
Type=oneshot
ExecStart=@BASE_BINDIR@/sh @BINDIR@/resin-proxy-config
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
