# meta-balena-common/recipes-core/systemd/disk-watchdog/disk-watchdog-failure.service
[Unit]
Description=Disk watchdog failure handler
DefaultDependencies=false

[Service]
Type=oneshot
ExecStart=/usr/bin/disk-watchdog-failure-handler
RemainAfterExit=no
TimeoutSec=30s
Restart=no
